# src/main/resources/upstart.conf
start on runlevel [2345]
stop on runlevel [!2345]

script
	chdir JARPATH
	exec java -jar JARNAME*.jar 
end script 